# preparation/roles/da.traefik_k8s/templates/traefik.service.j2
[Unit]
Description=Traefik

[Service]
Type=notify
ExecStart=/usr/local/bin/traefik --configFile=/etc/traefik.toml
Restart=always
WatchdogSec=1s

[Install]
WantedBy=multi-user.target